# docker/php/crontab
* * * * * echo "Hello world" >> /var/log/cron.log
* * * * * /bin/bash /etc/scheduler.sh
